# packing/postinst/postinst
#!/bin/sh
set -e
systemctl daemon-reload
systemctl enable control_program.service
exit 0
